# package/x11/xfree86/xfree86.cron
#!/bin/sh

echo "Running mkfontscale / mkfontdir ..."
for dir in /usr/X11R6/lib/X11/fonts/* ; do
	echo -n "$dir "
	mkfontscale $x
	mkfontdir $x
done ; unset dir
echo

if [ "`which fc-cache`" ] ; then
	echo "Running fc-cache ..."
	fc-cache -v
fi

