# package/x11/xfree86/xfree86.cron
#!/bin/sh

echo "Running mkfontdir / mkttfontdir ..."
for dir in /usr/X11R6/lib/X11/fonts/* ; do
	ttfiles="`ls $dir/*.ttf $dir/*.TTF 2> /dev/null`"
	if [ "$ttfiles" ] ; then
		( cd $dir
		mkttfontdir $x > /dev/null )
	fi
	mkfontdir $dir || true
	unset ttfiles
done ; unset dir

if [ "`which fc-cache`" ] ; then
	echo "Running fc-cache ..."
	fc-cache -v
fi

