# target/tlx/rootfs/etc/inittab
# TLX demo initrd, based on the busybox example

# Start an "askfirst" shell on the console (whatever that may be)
::askfirst:-/bin/sh
# Start an "askfirst" shell on /dev/tty2-4
vc/2::askfirst:-/bin/sh
vc/3::askfirst:-/bin/sh
vc/4::askfirst:-/bin/sh

# /sbin/getty invocations for selected ttys
vc/5::respawn:/sbin/getty 38400 vc/5
vc/6::respawn:/sbin/getty 38400 vc/6

# Example of how to put a getty on a serial line (for a terminal)
#::respawn:/sbin/getty -L ttyS0 9600 vt100
#::respawn:/sbin/getty -L ttyS1 9600 vt100
#
# Example how to put a getty on a modem line.
#::respawn:/sbin/getty 57600 ttyS2

# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
