# package/gnome2/gconf22/gconf.cron
#!/bin/sh

GNOME2_DIR="D_prefix"

for each in /etc/$GNOME2_DIR/gconf/schemas/*.schemas; do
	GCONF_CONFIG_SOURCE=xml::/etc/$GNOME2_DIR/gconf/gconf.xml.defaults gconftool-2 --makefile-install-rule $each
done
