# package/gnome2/scrollkeeper10/scrollkeeper.cron
#!/bin/sh

echo "Rebuilding scrollkeeper database"

if [ -x @GNOME_PREFIX@/bin/scrollkeeper-rebuilddb ]; then
	echo "GNOME-1 ...."
	@GNOME_PREFIX@/bin/scrollkeeper-rebuilddb
fi

if [ -x @GNOME2_PREFIX@/bin/scrollkeeper-rebuilddb ]; then
	echo "GNOME-2 ...."
	@GNOME2_PREFIX@/bin/scrollkeeper-rebuilddb
fi
