# package/gnome2/scrollkeeper10/scrollkeeper.cron
#!/bin/sh

GNOME2_DIR="D_prefix"
if [ -x /$GNOME2_DIR/bin/scrollkeeper-rebuilddb ]; then
	echo "Rebuilding scrollkeeper database"
	/$GNOME2_DIR/bin/scrollkeeper-rebuilddb
fi
