# package/base/hal/70-hal.rules

ENV{UDEVD_EVENT}=="1",	RUN+="/usr/libexec/hal.hotplug"

