# package/base/hal/70-hal.rules

# call the hal "hotplug" helper for device changes
ACTION=="add", SUBSYSTEM=="*", RUN+="/usr/libexec/hal.hotplug $env{SUBSYSTEM}"

