# package/base/prelink/prelink.cron
#!/bin/sh

echo "Running prelink utility ..."
prelink -av

