# package/filesystem/lvm2/lvm2-activate.service
[Unit]
Description="Activate lvm2 volumes"
DefaultDependencies=no
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=-/sbin/vgchange -ay
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
