# package/x11/xfree86-cvs/xfree86-cvs.cron
#!/bin/sh

echo "Running mkfontdir ..."
for i in /usr/X11R6/lib/X11/fonts/* ; do
	if test -f $i/fonts.dir ; then mkfontdir $i ; fi
done ; unset i

