# package/gnome24/scrollkeeper/scrollkeeper.cron
#!/bin/sh

GNOME24_DIR="D_prefix"
if [ -x /$GNOME24_DIR/bin/scrollkeeper-updatedb ]; then
	echo "Updating scrollkeeper database"
	/$GNOME24_DIR/bin/scrollkeeper-updatedb
fi
