# package/init/systemd/systemd.postinstall
#!/bin/bash
# --- T2-COPYRIGHT-BEGIN ---
# t2/package/*/systemd/systemd.postinstall
# Copyright (C) 2022 - 2026 The T2 SDE Project
# SPDX-License-Identifier: GPL-2.0
# --- T2-COPYRIGHT-END ---

echo "Running systemctl default preset-all"

systemctl preset-all
systemctl disable {cxl,ndctl}-monitor.service
false && systemctl daemon-reload

systemctl --global enable {pipewire,pipewire-pulse}.socket wireplumber.service
