# target/archivista/rootfs/etc/inittab
#!patch

Boot into X by default and run the archivista services.

--- inittab	2005-07-20 18:47:23.000000000 +0200
+++ inittab	2005-07-20 19:05:38.000000000 +0200
@@ -24,7 +24,7 @@
 #
 
 # The default runlevel.
-id:2:initdefault:
+id:5:initdefault:
 
 # /etc/rc.d executes the S and K scripts upon change
 # of runlevel.
@@ -82,3 +82,7 @@
 #
 X:5:respawn:/usr/X11/bin/startxdm 
 
+# The Archivista jobs
+av1:5:respawn:/usr/bin/perl /home/cvs/archivista/jobs/sane-daemon.pl
+av2:5:respawn:/usr/bin/perl /home/cvs/archivista/jobs/createpdf.pl
+
